# images/php/supervisord.conf
[supervisord]
nodaemon=true
[program:phpfpm]
command=php-fpm

[program:pm2]
command=pm2 startOrRestart /opt/www/Vanilla.Api.Ele/apps_ele.json