# images/node/supervisord.conf
[supervisord]
nodaemon=true

[program:code.aacc.in]
command=pm2 startOrRestart /opt/www/code.aacc.in/code.aacc.in.sh

[program:n.aacc.in]
command=pm2 startOrRestart /opt/www/n.aacc.in/note.aacc.in.sh